# files/etc/uci-defaults/999-QOL_config
uci -q batch << EOI

# Wireless default settings for radio0
set wireless.default_radio0.wnm_sleep_mode='1'
set wireless.default_radio0.ieee80211k='1'
set wireless.default_radio0.time_advertisement='2'
set wireless.default_radio0.bss_transition='1'
set wireless.default_radio0.proxy_arp='1'
set wireless.default_radio0.ocv='1'

# Wireless radio0 specific settings
set wireless.radio0.htmode='HE160'
set wireless.radio0.channel='auto'
set wireless.radio0.country='CA'

# Wireless default settings for radio1
set wireless.default_radio1.wnm_sleep_mode='1'
set wireless.default_radio1.ieee80211k='1'
set wireless.default_radio1.time_advertisement='2'
set wireless.default_radio1.bss_transition='1'
set wireless.default_radio1.proxy_arp='1'
set wireless.default_radio1.ocv='1'

# Wireless radio1 specific settings
set wireless.radio1.htmode='HE40'
set wireless.radio1.channel='auto'
set wireless.radio1.country='CA'

# Disable Firewall service entirely
set firewall.@zone[0].input='ACCEPT'
set firewall.@zone[0].forward='ACCEPT'
set firewall.@zone[0].output='ACCEPT'

# Network settings for a pure bridge
set network.lan.proto='static'
set network.lan.ipaddr='192.168.1.2' # Set to your desired AP IP
set network.lan.netmask='255.255.255.0'
set network.lan.gateway='192.168.1.1' # Point to your main router
set network.lan.dns='192.168.1.1'

EOI
