# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# These services must not run on a Dumb AP
for i in dnsmasq odhcpd; do
    if /etc/init.d/"$i" enabled; then
        /etc/init.d/"$i" disable
        /etc/init.d/"$i" stop
    fi
done

exit 0
